# debian/python3-citrix-adc-frrouting-sync.postinst
#!/bin/sh
set -e

# Automatically added by dh_python3:
if which py3compile >/dev/null 2>&1; then
	py3compile -p python3-citrix-adc-frrouting-sync 
fi
if which pypy3compile >/dev/null 2>&1; then
	pypy3compile -p python3-citrix-adc-frrouting-sync  || true
fi

# End automatically added section

systemctl daemon-reload
systemctl enable citrix-adc-frrouting-sync.service
